# gateway-software/systemd/gateway-manager.service
[Unit]
Description=SiteConnector Gateway Manager
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/local/bin
Environment=MONITORING_ENABLED=true
Environment=AUTO_START_TUNNEL=true
Environment=LOG_LEVEL=INFO
ExecStart=/usr/bin/python3 /usr/local/bin/gateway_manager.py monitor
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=siteconnector-gateway

# Security settings (relaxed for gateway functionality)
NoNewPrivileges=no
PrivateTmp=yes
ReadWritePaths=/etc/siteconnector /etc/wireguard-gateway /etc/wireguard /var/log

[Install]
WantedBy=multi-user.target