# gateway-software/systemd/network-scanner.service
[Unit]
Description=SiteConnector Network Scanner
After=network.target siteconnector-gateway.service
Wants=network.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/usr/local/bin
Environment=LOG_LEVEL=INFO
ExecStart=/usr/bin/python3 /usr/local/bin/network-scanner.py
StandardOutput=journal
StandardError=journal
SyslogIdentifier=siteconnector-scanner

# Security settings (relaxed for network scanning)
NoNewPrivileges=no
PrivateTmp=yes
ReadWritePaths=/etc/siteconnector /etc/wireguard-gateway /var/log

[Install]
WantedBy=multi-user.target