# gateway-os/services/gateway-firstboot.service
[Unit]
Description=Gateway First Boot Setup
Documentation=https://github.com/your-repo/gateway-os
DefaultDependencies=false
After=systemd-remount-fs.service
Before=systemd-sysusers.service sysinit.target shutdown.target
Conflicts=shutdown.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/gateway-firstboot
StandardOutput=journal+console
StandardError=journal+console

# Security settings
User=root
Group=root

# Ensure this only runs once
ConditionPathExists=!/var/lib/gateway/.firstboot-complete

[Install]
WantedBy=sysinit.target